# logstash.conf
[program:logstash]
command=/logstash-1.5.2/bin/logstash --verbose -f /logstash-1.5.2/logstash-config.conf
autostart=false
autorestart=false
